# tiny_runtime/systemd/tinyhat-runtime-gateway.service
[Unit]
Description=Tinyhat OpenClaw gateway
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/tinyhat/runtime.env
WorkingDirectory=/opt/tinyhat/current
ExecStart=/opt/tinyhat/current/bin/tinyhat-runtime gateway run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
